# etc/systemd/sing-box.service
[Unit]
Description=5gpn sing-box (SNI/QUIC transparent proxy)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/sing-box run -c /usr/local/etc/sing-box/config.json
Restart=on-failure
RestartSec=5
User=root
LimitNOFILE=65535
# Sandbox: leaf proxy — binds 80/443 tcp + 443 udp, root with no new privileges,
# read-only config, IPv4 (+unix) sockets only. Writes nothing (logs -> journald).
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_INET AF_UNIX
ReadOnlyPaths=/usr/local/etc/sing-box

[Install]
WantedBy=multi-user.target
